# com.hechuan.sales-system-demo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hechuan.sales-system-demo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>http.server</string>
    <string>--bind</string>
    <string>0.0.0.0</string>
    <string>8088</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/private/tmp/sales-system-demo-public</string>
  <key>StandardOutPath</key>
  <string>/Users/hechuan/Documents/Codex/2026-07-05/new-chat/outputs/sales-system-demo/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hechuan/Documents/Codex/2026-07-05/new-chat/outputs/sales-system-demo/server.err.log</string>
</dict>
</plist>
